# package/mesh-ac/files/etc/init.d/mesh-ac
#!/bin/sh /etc/rc.common

START=70
USE_PROCD=1

start_service() {
	mkdir -p /etc/mesh-ac/nodes /tmp/mesh-ac
	procd_open_instance
	procd_set_param command /bin/sh -c 'sleep infinity'
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger mesh_ac
}

reload_service() {
	if [ -x /usr/sbin/mesh-ac-apply-local ]; then
		/usr/sbin/mesh-ac-apply-local >/dev/null 2>&1 || true
	fi
	restart
}
