# package/mesh-ac/files/etc/uci-defaults/90-mesh-ac-enable
#!/bin/sh

/etc/init.d/mesh-ac enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd restart >/dev/null 2>&1 || true

exit 0
