# package/mesh-ac/files/etc/uci-defaults/90-mesh-ac-enable
#!/bin/sh

set_ath11k_nss_offload() {
	local value="$1"

	[ -f /etc/modules.d/ath11k ] || return 0
	if grep -q 'nss_offload=' /etc/modules.d/ath11k; then
		sed -i "s/nss_offload=[01]/nss_offload=$value/g" /etc/modules.d/ath11k
	else
		sed -i "s/^ath11k\\([[:space:]]\\|$\\)/ath11k nss_offload=$value\\1/" /etc/modules.d/ath11k
	fi
	echo "$value" > /sys/module/ath11k/parameters/nss_offload 2>/dev/null || true
}

uci -q get mesh_ac.main.enabled >/dev/null 2>&1 || uci set mesh_ac.main.enabled='1'
uci -q get mesh_ac.main.pairing_enabled >/dev/null 2>&1 || uci set mesh_ac.main.pairing_enabled='1'
uci -q get mesh_ac.main.local_member >/dev/null 2>&1 || uci set mesh_ac.main.local_member='0'
uci -q get mesh_ac.main.network_mode >/dev/null 2>&1 || uci set mesh_ac.main.network_mode='bridge'
uci commit mesh_ac >/dev/null 2>&1 || true

if [ "$(uci -q get mesh_ac.main.local_member 2>/dev/null)" = "1" ]; then
	set_ath11k_nss_offload 0
else
	set_ath11k_nss_offload 1
fi

/etc/init.d/mesh-ac enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd restart >/dev/null 2>&1 || true

# Advertise the AC over mDNS so AP agents can discover it without a hardcoded URL.
/etc/init.d/umdns enable >/dev/null 2>&1 || true
/etc/init.d/umdns restart >/dev/null 2>&1 || true

exit 0
